# services/manage_autohat_slaves.service
[Unit]
Description=AutoHAT Jenkins Slave Manager

[Service]
Type=simple
ExecStart=/etc/autohat/services/manage_autohat_slaves.sh
